# roles/logs_server/files/rsyslog-restart.service
[Unit]
Description=Restart rsyslog.service
Requires=rsyslog.service

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart rsyslog
